# salt/openssh/files/update-motd.d/00-figlet
#!/bin/sh

. /etc/update-motd.d/colors

printf "\n"$LIGHT_GREEN
figlet -f script "  "$(hostname -s)
printf $NONE
printf "\n"
